# src/treasurehunt.socket
[Unit]
Description=FAUST CTF Treasure Hunt

[Socket]
ListenStream=12321
Accept=true

[Install]
WantedBy=sockets.target
